# packages/rsync/rsync.install
post_upgrade() {
	systemctl -q is-enabled rsyncd.service rsyncd.socket
}
